# config/systemd/greetd.service.d/deepblack-vt-palette.conf
[Service]
ExecStartPre=/usr/local/bin/deepblack-apply-vt-palette
