# systemd/user/deepblack-mako.service
[Unit]
Description=deepBlack source-tokenized Mako notification daemon
Documentation=https://github.com/emersion/mako
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/Projects/deepBlack/scripts/deepblack-mako.sh
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
